# deploy/all-in-one/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; Infrastructure
[program:postgresql]
command=/app/postgres-start.sh
autostart=true
autorestart=true
priority=10
startsecs=3
startretries=10
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql.err

[program:redis]
command=redis-server --daemonize no --bind 127.0.0.1 --port 6379
autostart=true
autorestart=true
priority=20
startsecs=1
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.err

[program:nats]
command=nats-server -js -m 8222
autostart=true
autorestart=true
priority=20
startsecs=1
stdout_logfile=/var/log/supervisor/nats.log
stderr_logfile=/var/log/supervisor/nats.err

; Wait for PostgreSQL then run migrations
[program:migrations]
command=/app/wait-for-db.sh /bin/sh -c "cat /app/migrations/*.sql | psql -h 127.0.0.1 -U ggid -d ggid -v ON_ERROR_STOP=1 -q && echo 'Migrations complete'"

[program:seed]
command=/app/wait-for-db.sh /app/seed.sh
autostart=true
autorestart=false
priority=35
startsecs=0
stdout_logfile=/var/log/supervisor/seed.log
stderr_logfile=/var/log/supervisor/seed.err
autostart=true
autorestart=false
priority=30
startsecs=0
stdout_logfile=/var/log/supervisor/migrations.log
stderr_logfile=/var/log/supervisor/migrations.err

; GGID Services
[program:identity-server]
command=/app/wait-for-db.sh /usr/local/bin/identity-server -http-addr :8081 -grpc-addr :50051
autostart=true
autorestart=true
priority=40
startsecs=2
stdout_logfile=/var/log/supervisor/identity-server.log
stderr_logfile=/var/log/supervisor/identity-server.err

[program:auth-server]
command=/app/wait-for-db.sh /usr/local/bin/auth-server
autostart=true
autorestart=true
priority=40
startsecs=2
stdout_logfile=/var/log/supervisor/auth-server.log
stderr_logfile=/var/log/supervisor/auth-server.err

[program:oauth-server]
command=/app/wait-for-db.sh /usr/local/bin/oauth-server -addr :9005 -private-key /app/configs/rsa_private.pem -public-key /app/configs/rsa_public.pem
autostart=true
autorestart=true
priority=40
startsecs=2
stdout_logfile=/var/log/supervisor/oauth-server.log
stderr_logfile=/var/log/supervisor/oauth-server.err

[program:policy-server]
command=/app/wait-for-db.sh /usr/local/bin/policy-server
autostart=true
autorestart=true
priority=40
startsecs=2
stdout_logfile=/var/log/supervisor/policy-server.log
stderr_logfile=/var/log/supervisor/policy-server.err

[program:org-server]
command=/app/wait-for-db.sh /usr/local/bin/org-server
autostart=true
autorestart=true
priority=40
startsecs=2
stdout_logfile=/var/log/supervisor/org-server.log
stderr_logfile=/var/log/supervisor/org-server.err

[program:audit-server]
command=/app/wait-for-db.sh /usr/local/bin/audit-server
autostart=true
autorestart=true
priority=40
startsecs=2
stdout_logfile=/var/log/supervisor/audit-server.log
stderr_logfile=/var/log/supervisor/audit-server.err

[program:gateway-server]
command=/app/wait-for-db.sh /usr/local/bin/gateway-server
autostart=true
autorestart=true
priority=50
startsecs=2
stdout_logfile=/var/log/supervisor/gateway-server.log
stderr_logfile=/var/log/supervisor/gateway-server.err

; Console
[program:console]
command=/bin/sh -c "cd /app/console && exec node server.js"
autostart=true
autorestart=true
priority=60
startsecs=3
stdout_logfile=/var/log/supervisor/console.log
stderr_logfile=/var/log/supervisor/console.err
